# etc/cron.d/update_wix_dns_cronjob
*/5 * * * * root /config/scripts/update_wix_dns.sh
